# google_compute_engine_init/openrc/google-ip-forwarding-daemon
#!/sbin/openrc-run
# Copyright 2016 Google Inc. All Rights Reserved.
#
# Licensed under the Apache License, Version 2.0 (the "License");
# you may not use this file except in compliance with the License.
# You may obtain a copy of the License at
#
#     http://www.apache.org/licenses/LICENSE-2.0
#
#     Unless required by applicable law or agreed to in writing, software
#     distributed under the License is distributed on an "AS IS" BASIS,
#     WITHOUT WARRANTIES OR CONDITIONS OF ANY KIND, either express or implied.
#     See the License for the specific language governing permissions and
#     limitations under the License.

# Do NOT "set -e".

name=google-ip-forwarding-daemon
command=/usr/bin/google_ip_forwarding_daemon
pidfile=/var/run/$name.pid
description='Google Compute Engine IP Forwarding Daemon'

# Exit if the package is not installed.
[ -x "$command" ] || exit 0

depend()
{
  need net syslog
  after google_network_setup
}

#
# Function that starts the daemon/service.
#
start()
{
  ebegin "Starting $SVCNAME"
  start-stop-daemon \
    --background \
    --exec $command \
    --make-pidfile \
    --pidfile $pidfile \
    --quiet \
    --start > /dev/null
  eend $?
}

#
# Function that stops the daemon/service.
#
stop()
{
  ebegin "Stopping $SVCNAME"
  start-stop-daemon \
    --exec $command \
    --pidfile $pidfile \
    --quiet \
    --retry=TERM/30/KILL/5 \
    --stop > /dev/null
  ok=$?

  # Wait for children to finish too if this is a daemon that forks
  # and if the daemon is only ever run from this initscript.
  # If the above conditions are not satisfied then add some other code
  # that waits for the process to drop all resources that could be
  # needed by services started subsequently.  A last resort is to
  # sleep for some time.
  start-stop-daemon \
    --exec $command \
    --oknodo \
    --quiet \
    --retry=0/30/KILL/5 \
    --stop > /dev/null

  # Delete the pidfile when the daemon exits.
  rm -f $pidfile
  eend $?
}
